# systemd/claudebot@.service
[Unit]
Description=Claudebot Discord Bot (%i)
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/claudebot/%i.env
WorkingDirectory=/opt/claudebot
ExecStart=/usr/bin/python3 /opt/claudebot/scripts/run_bot.py --instance %i --env-file /etc/claudebot/%i.env
Restart=on-failure
RestartSec=30
KillMode=mixed
TimeoutStopSec=30
StandardOutput=append:/var/log/claudebot/%i.log
StandardError=append:/var/log/claudebot/%i.log

[Install]
WantedBy=multi-user.target
